# systemd/dann_guard.service
[Unit]
Description=DANN Guard - Pterodactyl Disk/Abuse Monitor
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/dann_guard
ExecStart=/root/dann_guard/dann_guard
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
